# frontpocket.service
[Unit]
Description=FrontPocket TTS Server
# Wait for the audio session to be available
After=pipewire.service pipewire-pulse.service pulseaudio.service
Wants=pipewire.service pipewire-pulse.service

[Service]
Type=simple
WorkingDirectory=%h/FrontPocket
EnvironmentFile=%h/.config/FrontPocket/environment
ExecStart=%h/FrontPocket/venv/bin/python3 %h/FrontPocket/frontpocket_server.py
Restart=on-failure
RestartSec=5

# Give the model time to load before systemd considers it timed out
TimeoutStartSec=120

# Logging goes to journald via stdout (matches our logging setup)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=frontpocket

[Install]
WantedBy=default.target
